# deploy/supervisor/ir4.conf
; IR4 Supervisor programs (DOC-20)
; Install: /etc/supervisor/conf.d/ir4.conf
; PHP-FPM is managed by systemd — not Supervisor.

[program:ir4-reverb]
command=/usr/bin/php /var/www/ir4/artisan reverb:start --host=127.0.0.1 --port=8080
directory=/var/www/ir4
user=www-data
autostart=true
autorestart=true
stopwaitsecs=10
stdout_logfile=/var/log/ir4/reverb.log
stderr_logfile=/var/log/ir4/reverb.err.log

[program:ir4-queue-default]
command=/usr/bin/php /var/www/ir4/artisan queue:work redis --queue=default --sleep=1 --tries=3 --max-time=3600
directory=/var/www/ir4
user=www-data
numprocs=2
process_name=%(program_name)s_%(process_num)02d
autostart=true
autorestart=true
stopwaitsecs=30
stdout_logfile=/var/log/ir4/queue-default.log
stderr_logfile=/var/log/ir4/queue-default.err.log

[program:ir4-queue-ingest]
command=/usr/bin/php /var/www/ir4/artisan queue:work redis --queue=ingest --sleep=1 --tries=3 --max-time=3600
directory=/var/www/ir4
user=www-data
numprocs=3
process_name=%(program_name)s_%(process_num)02d
autostart=true
autorestart=true
stopwaitsecs=30
stdout_logfile=/var/log/ir4/queue-ingest.log
stderr_logfile=/var/log/ir4/queue-ingest.err.log

[program:ir4-queue-reports]
command=/usr/bin/php /var/www/ir4/artisan queue:work redis --queue=reports --sleep=1 --tries=3 --max-time=3600
directory=/var/www/ir4
user=www-data
numprocs=1
process_name=%(program_name)s_%(process_num)02d
autostart=true
autorestart=true
stopwaitsecs=60
stdout_logfile=/var/log/ir4/queue-reports.log
stderr_logfile=/var/log/ir4/queue-reports.err.log

[program:ir4-scheduler]
command=/usr/bin/php /var/www/ir4/artisan schedule:work
directory=/var/www/ir4
user=www-data
autostart=true
autorestart=true
stopwaitsecs=10
stdout_logfile=/var/log/ir4/scheduler.log
stderr_logfile=/var/log/ir4/scheduler.err.log

[group:ir4]
programs=ir4-reverb,ir4-queue-default,ir4-queue-ingest,ir4-queue-reports,ir4-scheduler
